# config/includes.chroot/etc/systemd/system/display-manager.service
[Unit]
Description=Display Manager
[Service]
ExecStart=/usr/bin/sddm
[Install]
WantedBy=graphical.target
